# october/.infra/terraform/files/etc/profile.d/bash_colors.sh
# 30m # Gray
# 31m # Red
# 32m # Green
# 33m # Yellow
# 34m # Dark blue
# 35m # Pink
# 36m # Blue
# 37m # White

export PS1='\[\033[01;33m\]\u@\h\[\033[0;37m\] \w \$\[\033[00m\] '
# export PS1='\[\033[01;31m\]\u@\h\[\033[0;31m\] \w \$\[\033[00m\] '
